# etc/supervisor/conf.d/logstash.conf
[program:logstash]
command=/bin/sh /etc/init.d/logstash start
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/logstash/stdout.log
stderr_logfile=/var/log/logstash/stderr.log