# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/bin/sh /etc/init.d/nginx start
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log